# weather_info_cron/weather_everyday.cron
5 0 * * * python3 weather_data_everyday.py
* * * * * python3 weather.py
